# jellyfin-media-player/udev/99-hdmi-audio.rules
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/bin/su carlos -c '/home/carlos/hdmi-audio-fix.sh'"
